# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:xvnc]
command=/usr/bin/Xtigervnc :1 -geometry %(ENV_DESKTOP_GEOMETRY)s -depth 24 -rfbport 5901 -rfbauth /home/vm/.vnc/passwd -SecurityTypes VncAuth -localhost=0 -AlwaysShared=1 -AcceptPointerEvents=1 -AcceptKeyEvents=1
autostart=true
autorestart=true
priority=10
user=vm
environment=HOME="/home/vm",USER="vm",DISPLAY=":1",XDG_RUNTIME_DIR="%(ENV_VM_RUNTIME_DIR)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:desktop]
command=/usr/local/bin/start-desktop
autostart=true
autorestart=true
priority=20
user=vm
environment=HOME="/home/vm",USER="vm",DISPLAY=":1",XDG_RUNTIME_DIR="%(ENV_VM_RUNTIME_DIR)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:chrome]
command=/usr/local/bin/start-chrome
autostart=true
autorestart=true
priority=30
user=vm
environment=HOME="/home/vm",USER="vm",DISPLAY=":1",XDG_RUNTIME_DIR="%(ENV_VM_RUNTIME_DIR)s",CHROME_REMOTE_DEBUGGING_PORT="9223"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:cdp-proxy]
command=/usr/bin/socat TCP-LISTEN:9222,bind=0.0.0.0,reuseaddr,fork TCP:127.0.0.1:9223
autostart=true
autorestart=true
priority=35
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --listen 6080 --vnc 127.0.0.1:5901
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
